# deploy/zentinel-waf.service
[Unit]
Description=Zentinel WAF Agent - ModSecurity-based Web Application Firewall
Documentation=https://github.com/zentinelproxy/zentinel
After=network.target zentinel-proxy.service
Requires=zentinel-proxy.service
PartOf=zentinel-proxy.service

[Service]
Type=notify
NotifyAccess=main
ExecStart=/usr/local/bin/zentinel-waf-agent
ExecReload=/bin/kill -USR1 $MAINPID
Restart=on-failure
RestartSec=5s
TimeoutStartSec=30s
TimeoutStopSec=30s
KillMode=mixed
KillSignal=SIGTERM

# User and group
User=zentinel
Group=zentinel

# Environment
Environment="RUST_LOG=info,zentinel_waf_agent=debug"
Environment="WAF_CONFIG=/etc/zentinel/waf.yaml"
EnvironmentFile=-/etc/zentinel/waf.env

# Working directory
WorkingDirectory=/var/lib/zentinel/waf

# Runtime directory and mode
RuntimeDirectory=zentinel
RuntimeDirectoryMode=0755

# State directory
StateDirectory=zentinel/waf
StateDirectoryMode=0700

# Logs directory
LogsDirectory=zentinel-waf
LogsDirectoryMode=0750

# Configuration directory
ConfigurationDirectory=zentinel
ConfigurationDirectoryMode=0755

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
RemoveIPC=true
LockPersonality=true
ProtectClock=true
ProtectHostname=true
ProtectKernelLogs=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Allow network access for metrics endpoint
PrivateNetwork=false
IPAddressDeny=any
IPAddressAllow=localhost

# Resource limits
MemoryMax=2G
MemoryHigh=1G
CPUQuota=200%
TasksMax=4096
LimitNOFILE=65535
LimitNPROC=512

# Capabilities (none needed)
CapabilityBoundingSet=
AmbientCapabilities=

# Read-write paths for WAF
ReadWritePaths=/var/log/zentinel-waf
ReadWritePaths=/var/run/zentinel
ReadWritePaths=/var/lib/zentinel/waf

# Read-only paths for rules and libraries
ReadOnlyPaths=/usr/share/modsecurity-crs
ReadOnlyPaths=/usr/local/share/modsecurity-crs
ReadOnlyPaths=/opt/modsecurity-crs
ReadOnlyPaths=/etc/zentinel
ReadOnlyPaths=/usr/lib
ReadOnlyPaths=/usr/local/lib

# Process options
Nice=0
IOSchedulingClass=best-effort
IOSchedulingPriority=4
CPUSchedulingPolicy=other
CPUSchedulingPriority=0

# Watchdog
WatchdogSec=30s

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=zentinel-waf

[Install]
WantedBy=multi-user.target
