# deploy/zentinel-ratelimit-agent.service
[Unit]
Description=Zentinel Rate Limit Agent - Token Bucket Rate Limiting
Documentation=https://github.com/zentinelproxy/zentinel
After=network.target zentinel.service
PartOf=zentinel.service
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
ExecStart=/usr/local/bin/ratelimit-agent --socket /var/run/zentinel/ratelimit.sock
Restart=on-failure
RestartSec=5s
TimeoutStartSec=10
TimeoutStopSec=10
KillMode=mixed
KillSignal=SIGTERM

# User and permissions
User=zentinel
Group=zentinel
UMask=0077

# Environment
Environment="RUST_LOG=info"
Environment="RATELIMIT_AGENT_SOCKET=/var/run/zentinel/ratelimit.sock"
Environment="RATELIMIT_AGENT_DEFAULT_RPS=100"
Environment="RATELIMIT_AGENT_DEFAULT_BURST=200"
EnvironmentFile=-/etc/zentinel/agents/ratelimit.env

# Configuration file
Environment="RATELIMIT_AGENT_CONFIG=/etc/zentinel/agents/ratelimit.yaml"

# Working directory
WorkingDirectory=/var/lib/zentinel
RuntimeDirectory=zentinel
RuntimeDirectoryMode=0755
StateDirectory=zentinel/ratelimit
StateDirectoryMode=0700

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
RemoveIPC=true
LockPersonality=true
ProtectClock=true
ProtectHostname=true
ProtectKernelLogs=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
SystemCallArchitectures=native
PrivateNetwork=true
PrivateDevices=true

# Resource limits
LimitNOFILE=16384
LimitNPROC=512
LimitCORE=0
TasksMax=512
CPUWeight=75
MemoryMax=512M
MemoryHigh=400M

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=zentinel-ratelimit-agent

# Process management
OOMScoreAdjust=-200

[Install]
WantedBy=zentinel.service
