# deploy/zentinel-echo-agent.service
[Unit]
Description=Zentinel Echo Agent - Reference Implementation
Documentation=https://github.com/zentinelproxy/zentinel
After=network.target zentinel.service
PartOf=zentinel.service
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
ExecStart=/usr/local/bin/echo-agent --socket /var/run/zentinel/echo.sock
Restart=on-failure
RestartSec=5s
TimeoutStartSec=10
TimeoutStopSec=10
KillMode=mixed
KillSignal=SIGTERM

# User and permissions
User=zentinel
Group=zentinel
UMask=0077

# Environment
Environment="RUST_LOG=info"
Environment="ECHO_AGENT_SOCKET=/var/run/zentinel/echo.sock"
Environment="ECHO_AGENT_PREFIX=X-Echo-"
EnvironmentFile=-/etc/zentinel/agents/echo.env

# Working directory
WorkingDirectory=/var/lib/zentinel
RuntimeDirectory=zentinel
RuntimeDirectoryMode=0755

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
RemoveIPC=true
LockPersonality=true
ProtectClock=true
ProtectHostname=true
ProtectKernelLogs=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
SystemCallArchitectures=native
PrivateNetwork=true
PrivateDevices=true

# Resource limits
LimitNOFILE=8192
LimitNPROC=512
LimitCORE=0
TasksMax=256
CPUWeight=50
MemoryMax=256M
MemoryHigh=200M

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=zentinel-echo-agent

# Process management
OOMScoreAdjust=-100

[Install]
WantedBy=zentinel.service
