# deploy/zentinel.service
[Unit]
Description=Zentinel Security-First Reverse Proxy
Documentation=https://docs.zentinelproxy.io https://github.com/zentinelproxy/zentinel
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
ExecStart=/usr/local/bin/zentinel --config /etc/zentinel/zentinel.kdl
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -TERM $MAINPID
Restart=on-failure
RestartSec=5s
TimeoutStartSec=10
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

# User and permissions. The system user is declared in
# /usr/lib/sysusers.d/zentinel.conf and created via systemd-sysusers.
User=zentinel
Group=zentinel
UMask=0077

# Environment
Environment="RUST_LOG=info"
Environment="ZENTINEL_CONFIG=/etc/zentinel/zentinel.kdl"
EnvironmentFile=-/etc/zentinel/env

# Working directory and managed directories.
# systemd creates and chowns these to User= on each start.
WorkingDirectory=/var/lib/zentinel
RuntimeDirectory=zentinel
RuntimeDirectoryMode=0755
StateDirectory=zentinel
StateDirectoryMode=0700
LogsDirectory=zentinel
LogsDirectoryMode=0755
ConfigurationDirectory=zentinel
ConfigurationDirectoryMode=0755

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
RemoveIPC=true
LockPersonality=true
ProtectClock=true
ProtectHostname=true
ProtectKernelLogs=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
SystemCallArchitectures=native

# Resource limits
LimitNOFILE=65535
LimitNPROC=4096
LimitCORE=0
TasksMax=4096
CPUWeight=100
MemoryMax=2G
MemoryHigh=1500M

# Network capabilities.
# CAP_NET_BIND_SERVICE lets the proxy bind to ports < 1024 (e.g. 80, 443)
# without running as root. Edit listener addresses in /etc/zentinel/zentinel.kdl
# to use them.
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=zentinel

# OOM behavior: stay alive when the host is under memory pressure.
OOMScoreAdjust=-500

# Watchdog (optional, requires implementation in proxy)
# WatchdogSec=10s

[Install]
WantedBy=multi-user.target
Alias=zentinel-proxy.service
